# services/systemd/lair-keystore.service
[Unit]
Description=Lair Keystore Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
Environment="LAIR_DIR=/var/lib/holochain/data/lair"
Environment="LAIR_PASSWORD=secure-password"
ExecStart=/bin/bash -c 'printf "${LAIR_PASSWORD}" | /usr/local/bin/lair-keystore --lair-root ${LAIR_DIR} server --piped'
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings
ProtectSystem=full
ReadWritePaths=/var/lib/holochain
ReadWritePaths=/tmp
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target 